# supervisor_confs/lcd.conf
[program:lcd]
command=python -u /home/pi/code/rpi_daily_display/lcd/main.py
directory=/home/pi
stdout_logfile=/home/pi/logs/lcd.log
redirect_stderr=true
