# supervisor_confs/matrix.conf
[program:matrix]
command=python -u /home/pi/code/rpi_daily_display/lcd/matrix.py
directory=/home/pi
stdout_logfile=/home/pi/logs/matrix.log
redirect_stderr=true
