# supervisor_confs/pir.conf
[program:pir]
command=python -u /home/pi/code/rpi_daily_display/blink/pir.py
directory=/home/pi
stdout_logfile=/home/pi/pir.log
redirect_stderr=true
